# deployment/mc-dashboard-api.service
[Unit]
Description=MC Server Dashboard API Backend
Documentation=https://github.com/your-org/mc-server-dashboard-api
After=network.target
Wants=network.target

[Service]
Type=simple
User=$USER
Group=$USER
WorkingDirectory=/opt/mcs-dashboard/api
Environment=PYTHONPATH=/opt/mcs-dashboard/api
Environment=DATABASE_URL=sqlite:///./app.db
Environment=SECRET_KEY=your-secret-key-change-this-in-production
Environment=ALGORITHM=HS256
Environment=ACCESS_TOKEN_EXPIRE_MINUTES=30
ExecStart=/home/$USER/.local/bin/uv run fastapi dev --host 0.0.0.0 --port 8000
ExecReload=/bin/kill -USR1 $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mc-dashboard-api

# Health check
ExecStartPost=/bin/bash -c 'for i in {1..30}; do curl -s http://localhost:8000/docs >/dev/null && break || sleep 2; done'

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/mcs-dashboard/api
PrivateTmp=true

# Resource limits
LimitNOFILE=65536
TimeoutStartSec=120
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target